# scripts/crontab
@reboot sudo python /home/pi/shmile/scripts/button.py > /home/pi/button.log
